# install/etc/systemd/system/ntp.service.d/override.conf
[Service]
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitInterval=200 
